# deploy/launchd/com.alicia.knee-shoulder-stock.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alicia.knee-shoulder-stock</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/alicia/Desktop/#python/knee_shoulder_stock/deploy/launchd/run_daily_launchd.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/alicia/Desktop/#python/knee_shoulder_stock</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>16</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/alicia/Desktop/#python/knee_shoulder_stock/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/alicia/Desktop/#python/knee_shoulder_stock/logs/launchd.err.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
